# debian/prerm
#!/bin/bash
set -e

# Check if systemd is available
if [ -d /run/systemd/system ]; then
    systemctl stop super-training-collector.service || true
    systemctl disable super-training-collector.service || true
fi
