# debian/prerm
#!/bin/bash
set -e

# Stop and disable service before removal
systemctl stop super-training-collector.service || true
systemctl disable super-training-collector.service || true
